# root/etc/cron.d/auto-archiver
0 * * * * /usr/local/bin/auto-archiver.sh >> /var/log/cron.log 2>&1
